# root/etc/s6-overlay/s6-rc.d/svc-aria2/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
exec \
    s6-setuidgid $RUNUSER aria2c --conf-path=/config/aria2.conf
