# root/etc/s6-overlay/s6-rc.d/svc-darkhttpd/run
#!/usr/bin/with-contenv bash

exec \
    darkhttpd /app/ariang --port 8000
