# root/etc/services.d/darkhttpd/run
#!/usr/bin/with-contenv bash

exec \
    darkhttpd /ariang --port 80
